# templates/etc/rc.d/init.d/stash.erb
#! /bin/sh

### BEGIN INIT INFO
# Provides:          stash
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Initscript for Atlassian Stash
# Description:  Automatically start Atlassian Stash when the system starts up.
#               Provide commands for manually starting and stopping Stash.
### END INIT INFO

# Adapt the following lines to your configuration
# RUNUSER: The user to run Stash as.
RUNUSER="<%= scope.lookupvar('stash::params::user') %>"

# STASH_INSTALLDIR: The path to the Stash installation directory
STASH_INSTALLDIR="<%= scope.lookupvar('stash::params::webappdir') %>"

# STASH_HOME: Path to the Stash home directory
STASH_HOME="<%= scope.lookupvar('stash::params::homedir') %>"

# ==================================================================================
# ==================================================================================
# ==================================================================================

# PATH should only include /usr/* if it runs after the mountnfs.sh script
PATH=/sbin:/usr/sbin:/bin:/usr/bin
DESC="Atlassian Stash"
NAME=<%= scope.lookupvar('stash::params::product') %>
PIDFILE=$STASH_INSTALLDIR/work/catalina.pid
SCRIPTNAME=/etc/init.d/$NAME

# Read configuration variable file if it is present
[ -r /etc/default/$NAME ] && . /etc/default/$NAME

# Define LSB log_* functions.
# Depend on lsb-base (>= 3.0-6) to ensure that this file is present.
. /lib/lsb/init-functions


run_with_home() {
    if [ "$RUNUSER" != "$USER" ]; then
        su - "$RUNUSER" -c "export STASH_HOME=${STASH_HOME};${STASH_INSTALLDIR}/bin/$1"
    else
        export STASH_HOME=${STASH_HOME};${STASH_INSTALLDIR}/bin/$1
    fi
}

#
# Function that starts the daemon/service
#
do_start()
{
    run_with_home start-stash.sh
}

#
# Function that stops the daemon/service
#
do_stop()
{
    if [ -e $PIDFILE ]; then
      run_with_home stop-stash.sh
    else
      log_failure_msg "$NAME is not running."
    fi
}


case "$1" in
  start)
    [ "$VERBOSE" != no ] && log_warning_msg "Starting $DESC" "$NAME"
    do_start
    case "$?" in
        0|1) [ "$VERBOSE" != no ] && log_warning_msg 0 ;;
        2) [ "$VERBOSE" != no ] && log_warning_msg 1 ;;
    esac
    ;;
  stop)
    [ "$VERBOSE" != no ] && log_warning_msg "Stopping $DESC" "$NAME"
    do_stop
    case "$?" in
        0|1) [ "$VERBOSE" != no ] && log_warning_msg 0 ;;
        2) [ "$VERBOSE" != no ] && log_warning_msg 1 ;;
    esac
    ;;
  status)
       #if [ ! -e $PIDFILE ]; then
       #  log_failure_msg "$NAME is not running."
        # return 1
       #fi
       #status_of_proc -p $PIDFILE "" $NAME && exit 0 || exit $?
      ps -ef |grep $NAME |grep -v grep |awk '{ print $2 }' | <%= scope.lookupvar('stash::params::javahome') %>/bin/jps |grep -v Jps |grep -v grep > /dev/null
      RETVAL=$?
      if [ $RETVAL -eq 0 ];then
        echo "$NAME is running"
        exit $RETVAL
      else
        echo "$NAME is stopped"
        exit $RETVAL
      fi
      ;;
  restart|force-reload)
    #
    # If the "reload" option is implemented then remove the
    # 'force-reload' alias
    #
    log_warning_msg "Restarting $DESC" "$NAME"
    do_stop
    case "$?" in
      0|1)
        do_start
        case "$?" in
            0) log_warning_msg 0 ;;
            1) log_warning_msg 1 ;; # Old process is still running
            *) log_warning_msg 1 ;; # Failed to start
        esac
        ;;
      *)
        # Failed to stop
        log_warning_msg 1
        ;;
    esac
    ;;
  *)
    echo "Usage: $SCRIPTNAME {start|stop|status|restart|force-reload}" >&2
    exit 3
    ;;
esac